# systemd/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway - AGENTE_NAME_PLACEHOLDER
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/.openclaw
# .env e opcional: o '-' faz o systemd ignorar se o arquivo nao existir
EnvironmentFile=-/root/.openclaw/.env
ExecStart=/usr/local/bin/openclaw gateway run --force
Restart=always
RestartSec=10
TimeoutStopSec=15

# Security hardening
ProtectSystem=full
ProtectHome=false
NoNewPrivileges=true
PrivateTmp=true

# Resource limits (ajusta conforme a VPS)
LimitNOFILE=65536
MemoryHigh=3G
MemoryMax=4G

[Install]
WantedBy=multi-user.target
